# crontab
* * * * * /aws-scripts-mon/mon-put-instance-data.pl --agregated --auto-scaling --mem-util --mem-used --mem-avail --from-cron
